# deploy/colonyos-daemon.service
[Unit]
Description=ColonyOS Autonomous Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=colonyos
Group=colonyos
WorkingDirectory=/opt/colonyos/repo
ExecStart=/opt/colonyos/venv/bin/colonyos daemon --verbose
Restart=on-failure
RestartSec=30
WatchdogSec=120

# Security hardening
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=yes
ReadWritePaths=/opt/colonyos/repo /opt/colonyos/repo/.colonyos

# Environment — tokens loaded from drop-in file
EnvironmentFile=/opt/colonyos/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=colonyos

[Install]
WantedBy=multi-user.target
